# systemd/rcflow.service
[Unit]
Description=RCFlow Action Server
After=network.target

[Service]
Type=simple
User=rcflow
WorkingDirectory=/opt/rcflow
# Settings loaded from /opt/rcflow/settings.json by the application
ExecStart=/opt/rcflow/rcflow run
Restart=on-failure
RestartSec=5

# Allow rcflow to read git repos owned by other users (git >= 2.35.2 safe.directory check)
Environment="GIT_CONFIG_COUNT=1"
Environment="GIT_CONFIG_KEY_0=safe.directory"
Environment="GIT_CONFIG_VALUE_0=*"
# SSH key and other optional overrides (written by installer when an owner SSH key is found)
EnvironmentFile=-/opt/rcflow/env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=no
# Whole install dir is writable so settings.json (resolved by paths.get_data_dir)
# can be created/updated by the service user.
ReadWritePaths=/opt/rcflow
# HOME points at the install dir so XDG fallbacks land somewhere writable.
Environment="HOME=/opt/rcflow"
PrivateTmp=true

[Install]
WantedBy=multi-user.target
